# com.percy.orchestrator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.percy.orchestrator</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/jeralcooley/.openclaw/workspace/bin/percy-orchestrator.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/jeralcooley/.openclaw/workspace/.orchestrator.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jeralcooley/.openclaw/workspace/.orchestrator-error.log</string>
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
